# packages/orchestrator/internal/template/build/core/rootfs/files/serial-getty-autologin.service.tpl
{{- /*gotype:github.com/e2b-dev/infra/packages/orchestrator/internal/template/build/core/rootfs.templateModel*/ -}}
{{ .WriteFile "/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf" 0o644 }}

[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin root %I 115200,38400,9600 vt102
